# ap/rsyslog/rsyslog.upstart
description "rsyslog logging daemon"
version 8.4.2-1

oom score -1
start on filesystem
stop on runlevel [06]

respawn

exec /usr/bin/rsyslogd -n -i /var/run/rsyslog/rsyslog.pid
