# ap/nginx/nginx.upstart
description "nginx http server"
version 1.6.2

start on network-services
stop on runlevel [06]

expect fork
respawn

exec /usr/bin/nginx
