# ap/nginx/nginx.upstart
description "nginx http server"
version 1.6.2

start on network-services

expect fork
respawn

exec /usr/bin/nginx
