# ap/dcron/dcron.upstart
description "The dcron task scheduling daemon"

start on filesystem
respawn
expect fork

exec /usr/bin/crond -S -b
