# ap/man-db/man-db.upstart
description "Update the manpage index"
start on file FILE=/usr/man/
task

exec mandb
