# ap/openssh/openssh.upstart
description "openssh server"

start on network-services
stop on runlevel [06]

respawn

exec /usr/bin/sshd -D
