# ap/openssh/openssh.upstart
description "nginx http server"

start on network-services

respawn

exec /usr/bin/opensshd -D
