# ap/dbus/dbus.upstart
description "The D-Bus userspace message bus."
version 1.6.24

start on filesystem

respawn

pre-start script
/usr/bin/dbus-uuidgen --ensure
end script

exec /usr/bin/dbus-daemon --system --nofork
